# systemd/tpt-patrol.service
[Unit]
Description=TPT Patrol – File integrity scanner & process monitor
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/tpt-patrol --config /etc/tpt/patrol.toml
Restart=on-failure
RestartSec=5s
User=root
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tpt-patrol
# Patrol needs read access everywhere; root is simplest for v1

[Install]
WantedBy=multi-user.target
